# root/etc/s6-overlay/s6-rc.d/haproxy/run
#!/command/with-contenv sh
exec haproxy -p /tmp/haproxy.pid -f /etc/haproxy/haproxy.cfg -db
